# nimble-security-oauth2/config/jboss/jboss.cron.d
* * * * * root /etc/jboss-as/jboss-up.sh > /tmp/jboss.cron